# deploy/update-agent/bastion-updater.service
[Unit]
Description=Bastion Update Agent
After=network.target
# On relay VM, start after relay:
# After=bastion-relay.service

[Service]
Type=simple
User=bastion-updater
Group=bastion-updater
WorkingDirectory=/opt/bastion-updater
ExecStart=/usr/bin/node /opt/bastion-updater/dist/main.js
Restart=always
RestartSec=10

# Security hardening
# NoNewPrivileges=no because agent needs sudo for git/pnpm/systemctl.
# The sudoers whitelist restricts what it can actually do.
NoNewPrivileges=no
ProtectHome=true
PrivateTmp=true

# NOTE: ProtectSystem=strict requires ReadWritePaths for every path the agent
# touches. Since the build path varies per VM (/opt/bastion on relay,
# /opt/bastion-ai on AI VM), we omit ProtectSystem and rely on the sudoers
# whitelist as the security boundary. Deployers who want stricter sandboxing
# can add a drop-in override:
#   /etc/systemd/system/bastion-updater.service.d/paths.conf
#   [Service]
#   ProtectSystem=strict
#   ReadWritePaths=/opt/bastion /var/lib/bastion

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
